# api/supervisord.d/transSync_tDVV.conf
[program:transSync_tDVV_0]
command     = /bin/bash -c 'source "$0" && exec "$@"' /etc/profile php transSync.php --chainid=tDVV --den=1 --isex=0
directory   = /root/elf_task/cli/
user        = root
startsecs   = 3
autorestart = true

redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = /srv/log/transSync_tDVV.log
