# api/supervisord.d/updateIndex.conf
[program:updateIndex]
command     = /bin/bash -c 'source "$0" && exec "$@"' /etc/profile php updateIndex.php
directory   = /root/elf_task/cli/
user        = root
startsecs   = 3
autorestart = true

redirect_stderr         = true
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_logfile          = /srv/log/updateIndex.log